# overlay/etc/s6-overlay/s6-rc.d/redis/run
#!/bin/withenv sh

set -euo pipefail


#======================================================================================================================
# Run preflight checks.
#======================================================================================================================

nu -c "use bf-redis run ; run preflight"


#======================================================================================================================
# Run redis server
#======================================================================================================================

s6-setuidgid redis redis-server ${BF_REDIS_CONF}
